# services/openclaw-backup.service
[Unit]
Description=OpenClaw daily backup

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/openclaw-backup-daily
